# custom_policies/files_repo/enc-openldap.service
[Unit]
Description=External Node Classification via OpenLDAP
After=network.target

[Service]
Type=forking
ExecStart=/export/content/openldap/libexec/slapd -u openldap -g openldap -h "ldap:// ldaps://"
Restart=always

[Install]
WantedBy=default.target
